# templates/default/delayed_worker.conf.erb
start on started mongodb
stop on stopping mongodb

script
exec sudo -u <%= @username %> <%= @user_path %>/delayed_job.sh start
end script

post-stop script
exec sudo -u <%= @username %> <%= @user_path %>/delayed_job.sh stop
end script