# templates/default/regenerate-secrets.service.erb
[Unit]
Description=Regenerate <%= @service_names.join(', ') %> secret keys
<% @service_names.each do |service| %>
Before=<%= service %>.service
<% end %>

[Service]
Type=oneshot
<% @service_names.each do |service| %>
ExecStart=/bin/sh -c "/usr/bin/<%= service %> config:set SECRET_KEY_BASE=$(/usr/bin/<%= service %> run rake secret)"
<% end %>
ExecStartPost=/bin/systemctl disable regenerate-secrets.service

[Install]
WantedBy=multi-user.target
