# templates/default/regenerate-secrets.service.erb
[Unit]
Description=Regenerate Cypress and CVU secret keys

[Service]
Type=oneshot
<% @secrets_paths.each do |secret_path| %>
ExecStart=/opt/regenerate-secrets.sh <%= secret_path %>
<% end %>
ExecStartPost=/bin/systemctl disable regenerate-secrets.service

[Install]
WantedBy=multi-user.target